# systemd/.config/systemd/user/mega_sync.service
[Unit]
Description=MegaSync Daemon

[Service]
ExecStart=%h/bin/mega_sync.sh
Restart=always

[Install]
WantedBy=xsession.target
